# src/xorif-startup/xorif-startup.service
[Unit]
Description=xorif-startup

[Service]
Type=simple
ExecStart=/usr/bin/xorif-startup.sh
StandardOutput=null

[Install]
WantedBy=multi-user.target